# voxy/core/infra/workstation/autoshutdown.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/autoshutdown.sh

[Install]
WantedBy=default.target
